# examples/launchagents/com.example.brew-autoupdate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<!-- Updates and upgrades Homebrew packages every day, cleaning up old versions. -->
<dict>
    <key>Label</key>
    <string>com.example.brew-autoupdate</string>

    <!-- Replace /Users/USERNAME below with your own home path (run `echo $HOME`). -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/USERNAME/.local/bin/brew-autoupdate.sh</string>
    </array>

    <!-- Run every day at 08:00 (no Weekday key => daily) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/USERNAME/Library/Logs/brew-autoupdate.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/USERNAME/Library/Logs/brew-autoupdate.log</string>
</dict>
</plist>
